# rag_server.service
[Unit]
Description=RAG System FastAPI Service
After=network.target

[Service]
# Port 80 requires root privileges
User=root
Group=root

# Change this to the actual directory on your server
WorkingDirectory=/home/asuengai/RAG/RAG

# Use the full path to the virtual environment's python
ExecStart=/home/asuengai/RAG/RAG/venv/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 80

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
